# rc.local
#!/bin/bash

for V in ${!CRON_D_*}
do
  F=/etc/cron.d/$(echo ${V#CRON_D_} | tr '[:upper:]' '[:lower:]')
  echo "${!V}" > $F
  chmod +x $F
done
