# rc.local
#!/bin/sh
find /etc/container_environment/CRON_D_* -print0 | xargs -0 -r ln -s -t /etc/cron.d/
